# trading_suite/systemd/trading-summary.service
[Unit]
Description=Trading Suite Unified Summary Collector
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/trading-bots
Environment=TRADING_PYTHON=/opt/trading-bots/.venv/bin/python
ExecStart=/opt/trading-bots/trading_suite/bin/run_summary_loop.sh
Restart=always
RestartSec=10
StandardOutput=append:/opt/trading-bots/trading_suite/runtime/summary/stdout.log
StandardError=append:/opt/trading-bots/trading_suite/runtime/summary/stderr.log

[Install]
WantedBy=multi-user.target
